# mqtt-listener.service
[Unit]
Description=Sistema de Presença - MQTT Listener
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/ser/projetos/sistema-presenca
ExecStart=/usr/bin/php /home/ser/projetos/sistema-presenca/artisan mqtt:listen
Restart=always
RestartSec=10
StandardOutput=append:/home/ser/projetos/sistema-presenca/storage/logs/mqtt-listener.log
StandardError=append:/home/ser/projetos/sistema-presenca/storage/logs/mqtt-listener-error.log

[Install]
WantedBy=multi-user.target